# systemd/fancontrol-hddtemp.service
[Unit]
Description=fancontrol-hddtemp

[Service]
Type=simple
ExecStart=/root/fancontrol-hddtemp

[Install]
WantedBy=multi-user.target
